# deploy/btbz-cms.service
# systemd 유닛 — /etc/systemd/system/btbz-cms.service (2026-07-20 서버 적용됨)
# 서버: 211.110.140.172 (btbz.ai · modora.btbz.ai · mail 등 공용) — btbz 사용자 + nvm node20 운영
[Unit]
Description=BTBZ CMS backend (NestJS)
After=network.target

[Service]
Type=simple
User=btbz
Group=btbz
WorkingDirectory=/home/btbz/btbz-cms
# .env: PORT=3110 / DB_PATH=/home/btbz/btbz-cms-data/... / JWT_SECRET / SMTP_* / SEED_ADMIN_* (600 권한)
EnvironmentFile=/home/btbz/btbz-cms/.env
ExecStart=/home/btbz/.nvm/versions/node/v20.20.2/bin/node dist/main.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
